# config/com.robinhood-ai-trader.shadow-worker-canary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.robinhood-ai-trader.shadow-worker-canary</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
    <string>/Users/ge/Documents/AI trading agent/robinhood-ai-trader/scripts/launchd_shadow_worker.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ge/Documents/AI trading agent/robinhood-ai-trader</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/ge/.local/bin:/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>ROBINHOOD_SHADOW_CANARY</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/ge/Documents/AI trading agent/robinhood-ai-trader/logs/launchd-canary.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ge/Documents/AI trading agent/robinhood-ai-trader/logs/launchd-canary.stderr.log</string>
</dict>
</plist>
